# etc/systemd/system/straylight-entropy.service
[Unit]
Description=StrayLight Hardware Entropy Daemon
Documentation=https://straylight.dev/docs/entropy
DefaultDependencies=no
After=sysinit.target local-fs.target

[Service]
Type=notify
ExecStart=/usr/bin/straylight-entropy
Restart=on-failure
RestartSec=1s
WatchdogSec=120s
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/straylight

[Install]
WantedBy=multi-user.target
