# etc/systemd/system/straylight-scheduler.service
[Unit]
Description=StrayLight Userspace Task Scheduler
After=straylight-registry.service
Requires=straylight-registry.service

[Service]
Type=simple
ExecStart=/usr/bin/straylight-scheduler
Restart=on-failure
RestartSec=2s
AmbientCapabilities=CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
