# etc/systemd/system/straylight-bus.service
[Unit]
Description=StrayLight D-Bus Session Manager
Documentation=https://straylight.dev/docs/bus
After=dbus.service
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/straylight-bus
Restart=on-failure
RestartSec=2s
User=root
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
